# deploy/systemd-sync.timer
# /etc/systemd/system/ctyrlistkoteka-sync.timer
#
# Enable with:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ctyrlistkoteka-sync.timer
#
# Inspect schedule / last run:
#   systemctl list-timers | grep ctyrlistkoteka
#   journalctl -u ctyrlistkoteka-sync.service -n 100
[Unit]
Description=Twice-daily Čtyřlístkotéka data sync
Requires=ctyrlistkoteka-sync.service

[Timer]
# Twice a day at 06:00 and 18:00 local — picks up rsync'd uploads
# within at most 12 hours without ever overlapping a peak-traffic slot.
OnCalendar=*-*-* 06,18:00:00
# Make up a missed run if the server was down at the scheduled time.
Persistent=true
# Spread by up to 10 min so multiple hosts wouldn't fire simultaneously.
RandomizedDelaySec=10m
Unit=ctyrlistkoteka-sync.service

[Install]
WantedBy=timers.target
